# supervisord.conf
[supervisord]
nodaemon=true

[program:cron]
priority=2
command = cron -f
autorestart=true

[program:redis]
priority=5
command=redis-server /etc/redis/redis.conf
stdout_logfile=/var/log/supervisor/redis.log
stderr_logfile=/var/log/supervisor/redis.log
autorestart=true

[program:postgresql]
command=/usr/lib/postgresql/9.5/bin/postgres -D /var/lib/postgresql/9.5/main/ -c config_file=/etc/postgresql/9.5/main/postgresql.conf
stdout_logfile=/var/log/supervisor/postgres.log
stderr_logfile=/var/log/supervisor/postgres.log
user=postgres
autorestart=true

[program:apache2]
priority=10
command=apache2ctl -DFOREGROUND
stdout_logfile=/var/log/supervisor/apache2.log
stderr_logfile=/var/log/supervisor/apache2.log
autorestart=true
